# arcane-gitops.service
[Unit]
Description=Docker Compose Git Sync and Arcane Deploy
After=network.target

[Service]
Type=oneshot
User=root
Group=root

# Environment file with configuration
EnvironmentFile=/etc/arcane-gitops/config.env

# Execute the sync binary
ExecStart=/usr/local/bin/arcane-gitops

# Logging
StandardOutput=journal
StandardError=journal
SyslogIdentifier=arcane-gitops

# Security hardening
NoNewPrivileges=true
PrivateTmp=true

[Install]
WantedBy=multi-user.target
